# etc/crontab
*/2 *   * * *   root    /usr/bin/python3 /home/section5_cron/scripts/myscript.py
